# src/settings/system/snap-snapd-21467.mount
[Unit]
Description=Mount unit for snapd, revision 21467
After=snapd.mounts-pre.target
Before=snapd.mounts.target
Before=local-fs.target

[Mount]
What=/var/lib/snapd/snaps/snapd_21467.snap
Where=/snap/snapd/21467
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
